# email/msmtp_runqueue.service
[Unit]
Description=run msmtp-runqueu.sh

[Service]
Type=oneshot
ExecStart=%h/config/email/msmtp_runqueue.sh
